# setup/manager_da.service
[Unit]
Description= Disk Analyzer manager startup and shutdown helper

#I will have to modify th onstart.sh path from setup to the working directory

[Service]
ExecStart=/home/characterme/Disk-Analyzer/manager/manager
KillMode=process
StandardOutput=journal+console
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
